# broker.go/systemd/monstermq-edge.service
[Unit]
Description=MonsterMQ Edge Broker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=monstermq
Group=monstermq
ExecStart=/usr/local/bin/monstermq-edge -config /etc/monstermq/config.yaml
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
